# systemd-units/mw-generateSitemap.timer
[Unit]
Description=Run Generate Sitemap Daily at 4:00 AM

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
